# infra/consul.service
[Unit]
Description=Consul
[Service]
ExecStart=consul agent -dev -client=0.0.0.0 -data-dir=/opt/consul
[Install]
WantedBy=multi-user.target