# infra/nomad.service
[Unit]
Description=Nomad
[Service]
ExecStart=nomad agent -config=/etc/nomad.d/nomad.hcl
[Install]
WantedBy=multi-user.target